# launchd/com.akimare.bot-yumekano-coe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.akimare.bot-yumekano-coe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/akimare/discord-bots/.venv/bin/python</string>
        <string>/Users/akimare/discord-bots/bot.py</string>
        <string>yumekano-coe</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/akimare/discord-bots</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/akimare/discord-bots/logs/yumekano-coe.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/akimare/discord-bots/logs/yumekano-coe.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
